# install/launchd-templates/zoekt-webserver.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.keisei.zoekt-webserver</string>

    <key>ProgramArguments</key>
    <array>
        <string>${BREW}/bin/zoekt-webserver</string>
        <string>-listen</string>
        <string>127.0.0.1:6070</string>
        <string>-index</string>
        <string>${DATA}/index</string>
    </array>

    <key>WorkingDirectory</key>
    <string>${DATA}</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>${HOME}</string>
        <key>USER</key>
        <string>${USER}</string>
        <key>PATH</key>
        <string>${BREW}/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>${LOGS}/webserver.out.log</string>

    <key>StandardErrorPath</key>
    <string>${LOGS}/webserver.err.log</string>

    <key>ProcessType</key>
    <string>Background</string>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>4096</integer>
    </dict>

    <key>HardResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>8192</integer>
    </dict>
</dict>
</plist>
